# @reboot-ec2/cold-db/etc/rc.local
#!/bin/sh -e
su ubuntu -c 'sh /home/ubuntu/cold_db-rc.sh >> /home/ubuntu/rc.local.log 2>&1 &'
exit 0