# @reboot-ec2/cache/etc/rc.local
#!/bin/sh -e
su ubuntu -c 'sh /home/ubuntu/cache-rc.sh >> /home/ubuntu/rc.local.log 2>&1 &'
exit 0